# supervisor.conf.example
# Template para Supervisor - Queue Worker
# Copiar a: /etc/supervisor/conf.d/offside-queue.conf
# Luego: supervisorctl reread && supervisorctl update

[program:offside-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/offside-app/artisan queue:work redis --sleep=3 --tries=3 --timeout=900 --max-jobs=1000
autostart=true
autorestart=true
stopasgroup=true
stopwaitsecs=3600
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/offside-queue-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
user=www-data
environment=PATH="/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html/offside-app"

[program:offside-schedule]
process_name=%(program_name)s
command=php /var/www/html/offside-app/artisan schedule:work
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/offside-schedule.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
user=www-data
stopasgroup=true
stopwaitsecs=600
environment=PATH="/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html/offside-app"

[group:offside-workers]
programs=offside-queue-worker,offside-schedule
priority=999
